# deploy/tbuddy.service
[Unit]
Description=T.Buddy translation relay (Gunicorn)
After=network.target

[Service]
Type=simple
User=tbuddy
Group=tbuddy
WorkingDirectory=/opt/tbuddy
# Use an EnvironmentFile to keep secrets out of systemd unit; create /etc/tbuddy/env
EnvironmentFile=/etc/tbuddy/env
ExecStart=/opt/tbuddy/venv/bin/gunicorn --bind 127.0.0.1:8080 --workers 1 app:app
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
